# contrib/caic.service
# caic systemd user service
# Install: mkdir -p ~/.config/systemd/user && cp contrib/caic.service ~/.config/systemd/user/
#          mkdir -p ~/.config/caic && cp contrib/caic.env ~/.config/caic/caic.env
#          Edit ~/.config/caic/caic.env to set CAIC_HTTP, CAIC_ROOT, and API keys.
# Enable:  systemctl --user daemon-reload && systemctl --user enable --now caic
# Logs:    journalctl --user -u caic -f

[Unit]
Description=Coding Agents in Containers
Documentation=https://github.com/maruel/caic
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-%h/.config/caic/caic.env
Environment=PATH=%h/go/bin:%h/src/md:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/go/bin/caic
WorkingDirectory=%h/src
Restart=always
RestartSec=1s

# --- Security Hardening ---

# Filesystem isolation
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/src
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectProc=invisible
ProcSubset=pid

# Capability restrictions
NoNewPrivileges=true
AmbientCapabilities=

# Network restrictions (allow only IPv4/IPv6/Unix sockets)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# System call filtering
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete

# Memory/execution restrictions
MemoryDenyWriteExecute=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true

# Misc
UMask=0077
KeyringMode=private
RemoveIPC=true

[Install]
WantedBy=default.target
